# scute.cron
*  *  *  *  *  python /usr/local/bin/updateip.py
