# scute.cron
*  *  *  *  *  python /usr/local/sbin/updateip.py
